# dist/xdg-desktop-portal-shojiwm.service
[Unit]
Description=Portal service (ShojiWM implementation)
PartOf=graphical-session.target
After=graphical-session.target

[Service]
Type=dbus
BusName=org.freedesktop.impl.portal.desktop.shojiwm
ExecStart=/usr/bin/xdg-desktop-portal-shojiwm
Restart=always
RestartSec=500ms
TimeoutStopSec=10
